# supervisor.conf
[supervisord]
logfile=/dev/stdout 
logfile_maxbytes=0  
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[unix_http_server]
file=/tmp/supervisor.sock

#https://stackoverflow.com/questions/25195880/supervisord-redirect-process-stdout-to-console
[program:indexer]
command=/srv/www/build-event-store
logfile=/dev/stdout 
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:api]
command=/srv/www/carbonmeter-api
logfile=/dev/stdout 
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
